# ugw/feeds_arc/arcadyan-utility/files/base-files/etc/profile.d/arc_profile_path.sh
#!/bin/sh
# We hiject switch_cli. So when shell calls switch_cli,
# we want shell use /bin/switch_cli instead of /usr/bin/switch_cli.
# Therefore, we put /bin in the begining of searching path.
export PATH=/bin:$PATH
